# etc/local.d/elogind.start
sudo bash -c 'echo 11680 > /sys/class/backlight/intel_backlight/brightness'
