# etc/local.d/elogind.start
sudo bash -c 'echo 11970 > /sys/class/backlight/intel_backlight/brightness'
polychromatic-cli -e Wave
